# backend/trip-planner-mcp.service
[Unit]
Description=Trip Planner MCP Server
After=network.target

[Service]
Type=simple
User=work
WorkingDirectory=/home/work/project/tripe-planner
Environment=PYTHONIOENCODING=utf-8
ExecStart=/usr/bin/python3 /home/work/project/tripe-planner/backend/mcp_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
